# distribution/deb/src/main/packaging/scripts/prerm
${packaging.scripts.header}

if [ -x "/etc/init.d/cassandra" ] || [ -e "/etc/init/cassandra.conf" ]; then
	invoke-rc.d cassandra stop || exit $?
fi

${packaging.scripts.footer}

